# bomba/control-bomba.service.in
[Unit]
Description=Water Bomb Pressure Monitor
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/python3 @SCRIPT_PATH@

[Install]
WantedBy=multi-user.target